# etc/supervisor.d/caddy.ini
[program:caddy]
command=caddy -conf /etc/caddy/caddy.conf
