# supervisor.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
autostart=True
autorestart=True